# contrib/com.vikram.team.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vikram.team</string>
    <key>ProgramArguments</key>
    <array>
        <string>__VIKRAM_WRAPPER__</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>VIKRAM_HOME</key>
        <string>__VIKRAM_HOME__</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>__VIKRAM_HOME__</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>__VIKRAM_HOME__/logs/gateway.log</string>
    <key>StandardErrorPath</key>
    <string>__VIKRAM_HOME__/logs/gateway.err</string>
    <key>ProcessType</key>
    <string>Background</string>
    <key>Nice</key>
    <integer>5</integer>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
